# servers/robot_controller_backend/network/omega-network-watchdog.service
[Unit]
Description=Omega Network Watchdog
Documentation=file:///home/omega1/Omega-Code/servers/robot_controller_backend/network/OMEGA_NETTOGGLE.md
After=network.target NetworkManager.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/omega-network-watchdog.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security
User=root
NoNewPrivileges=true

# Resource limits
MemoryLimit=50M
CPUQuota=10%

[Install]
WantedBy=multi-user.target

